# smlight_slzb07_programmer/rootfs/etc/s6-overlay/s6-rc.d/smlight-slzb-07-programmer/up
/etc/s6-overlay/scripts/smlight-slzb-07-programmer.sh